# overlay/etc/cron.d/observium
# Run discovery and poller wrappers
33  */6 * * * root /opt/observium/observium-wrapper discovery >> /dev/null 2>&1
*/5  *  * * * root /opt/observium/observium-wrapper discovery -w 1 --host new >> /dev/null 2>&1
*/5  *  * * * root /opt/observium/observium-wrapper poller >> /dev/null 2>&1

# Run housekeeping script daily for syslog, eventlog and alert log
13   5  * * * root /opt/observium/housekeeping.php -ysel >> /dev/null 2>&1
# Run housekeeping script daily for rrds, ports, orphaned entries in the database and performance data
47   4  * * * root /opt/observium/housekeeping.php -yrpb >> /dev/null 2>&1
